# systemd/user/conky.service
[Unit]
Description=Conky

[Service]
ExecStart=conky --daemonize --pause=5 -c ~/.config/conky/config

[Install]
WantedBy=default.target
